# imac/com.airtrackr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    AirTrackr — Tracker LaunchAgent (BACKUP ONLY)

    NOTE: The tracker requires Terminal.app for Accessibility permissions.
    Primary startup is via AirTrackr.command (Login Item).
    This plist is kept as a backup but should NOT be loaded simultaneously.
    Use com.airtrackr.api.plist for the API server instead.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airtrackr</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/evelyn/Repos/airtrackr/start_all_imac.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/evelyn/Repos/airtrackr/logs/tracker_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/evelyn/Repos/airtrackr/logs/tracker_launchd_stderr.log</string>
</dict>
</plist>
