# imac/com.airtrackr.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.airtrackr.cleanup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/evelyn/Repos/airtrackr/imac/cleanup_processes.sh</string>
    </array>

    <!-- Run 60 seconds after login (give time for apps to launch first) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>60</integer>

    <!-- Only run once at startup, then disable interval -->
    <key>LaunchOnlyOnce</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/evelyn/Repos/airtrackr/logs/cleanup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/evelyn/Repos/airtrackr/logs/cleanup.log</string>
</dict>
</plist>
